# misc/crontabs.cron
*/1 * * * * /home/m1ch3l/m1ch3l/misc/sync.bash
29 0 * * * /home/m1ch3l/m1ch3l/crons/il-est-minuit-trente.sh
23 8 * * * /home/m1ch3l/m1ch3l/crons/bonjour-du-matin.sh
15 12 * * * /home/m1ch3l/m1ch3l/crons/il-faut-manger.sh
